# webradio_0.1-1/etc/systemd/system/spotify.service
[Unit]
Description=spotify
After=network.target

[Service]
WorkingDirectory=/home/pi
Home=/home/pi
ExecStart=/usr/bin/java -jar /opt/librespot/librespot-player.jar
#Restart=always
User=pi

[Install]
WantedBy=multi-user.target